# cubeos-watchdog.service
[Unit]
Description=CubeOS Self-Healing Watchdog

[Service]
Type=oneshot
ExecStart=/cubeos/scripts/watchdog-health.sh
